# openrc/rflinkproxy
#!/sbin/openrc-run

command="/usr/local/bin/rflinkproxy"
if [[ -z "${HOST}" ]];
then
	# Serial mode.
	command_args="--listenport=${LISTEN_PORT} --port=${CONNECT_PORT} --baud=${BAUD_RATE}"
else
	# TCP mode.
    command_args="--listenport=${LISTEN_PORT} --port=${CONNECT_PORT} --host=${HOST}"
fi
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"
